# nexus-oss/DEBIAN/postinst
#!/bin/sh
# postinst script for mongodb
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
	# create a nexus group and user
        if ! grep -q nexus /etc/passwd; then
		adduser --system --no-create-home nexus
		addgroup --system nexus
		adduser nexus nexus
        fi

	# create db -- note: this should agree with dbpath in nexus.conf
	mkdir -p /var/lib/nexus
	chown -R nexus:nexus /var/lib/nexus

	chown -R nexus:nexus /var/lib/sonatype-work

	mkdir -p /var/run/nexus
	chown -R nexus:nexus /var/run/nexus

	mkdir -p /var/tmp/nexus
	chown -R nexus:nexus /var/tmp/nexus

	# create logdir -- note: this should agree with logpath in nexus.conf
	mkdir -p /var/log/nexus
	chown -R nexus:nexus /var/log/nexus

	chown -R nexus:nexus /usr/local/nexus*

	a2enmod proxy
	a2enmod proxy_http
	service apache2 restart

	 update-rc.d nexus start 20 2 3 4 5 . stop 0 1 6 .
	 /etc/init.d/nexus start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0